# upstart/address-book-updater.conf.in
description "address-book-update"
author "Renato Araujo Oliveira Filho <renato.filho@canonical.com>"

start on started address-book-service and started msyncd and started sync-monitor
stop on runlevel [06]

exec @CMAKE_INSTALL_FULL_LIBEXECDIR@/address-book-updater --sync
